# docker/all-in-one/supervisor/supervisord.conf
[program:nodejs]
command=yarn start
directory=/app/frontend
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
redirect_stderr=true
stdout_logfile_maxbytes = 0
stderr_logfile_maxbytes = 0
environment=NODE_ENV="production",VITE_API_URL_CLIENT="http://app.ticketz.com.br",VITE_API_URL_SERVER="http://app.ticketz.com.br/api",VITE_FRONTEND_URL="http://app.ticketz.com.br",VITE_STRIPE_PUBLISHABLE_KEY="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"
